# packaging/macos/com.mythodikal.mythd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist for the Mythodikal Anti-Virus user-level daemon
  (TASK-080 / TASK-254 — Phase 9 Wave 1 + Wave 2).

  Installed by the .dmg first-run flow at
  `~/Library/LaunchAgents/com.mythodikal.mythd.plist` and loaded with
  `launchctl bootstrap gui/<uid>`. Per `docs/prd.md` § 1.5.3 the plist
  is unsigned — the macOS first-launch prompt is documented in
  `docs/macos-first-run.md`.

  The plist is per-user (`LaunchAgents`, not `LaunchDaemons`) so no
  sudo install is required. KeepAlive bounces the daemon if it
  crashes; the heartbeat file under
  `~/Library/Application Support/Mythodikal/heartbeat.json` is what
  the Real-time page surfaces as a green/amber/red chip (TASK-254).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mythodikal.mythd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mythd-macos</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>1</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <!--
      StandardOutPath / StandardErrorPath intentionally omitted.
      launchd does not expand `~` or `$HOME` in plist path keys, so
      routing logs through a per-user file would require the .dmg
      first-run installer to rewrite the plist with the resolved
      $HOME at install time. Until that installer pass lands,
      letting launchd route stdout / stderr to the system asl log
      (visible in Console.app under the daemon's Label) keeps logs
      off world-readable `/tmp/` — a non-root attacker on the same
      Mac could otherwise watch the file for daemon internals
      (security-review follow-up, Phase 9 Wave 2 closeout).
    -->
</dict>
</plist>
